# Lora-Bot/dpkg/postinst
#!/bin/bash

systemctl enable lorabot
systemctl daemon-reload

if [ -f /tmp/lorabot_service_runner ]; then
    service lorabot start
    rm /tmp/lorabot_service_runner
fi