# Lora-Bot/dpkg/postinst
#!/bin/bash

systemctl enable lorabot
systemctl daemon-reload
if [[ $(systemctl is-active lorabot || true) == "active" ]]
then
        service lorabot restart
fi
